# cfinit/a/service/Init/run
#!/bin/dash
trap "trap - TERM && kill -- -$$" INT TERM EXIT

while date; do
  sleep 5
done

