# cfinit/a/service/request-proxy/run
#!/bin/bash
trap "trap - TERM && kill -- -$$" INT TERM EXIT
